# data/pay-service.conf.in
description "Service to process payments"

start on desktop-start
stop on desktop-end

exec @CMAKE_INSTALL_LIBEXECDIR@/pay-service
