# data/ubuntu-pay-service-trust-stored.conf
description "Pay Service Trust Store Daemon"

# Try to start this once pay-service finishes starting
start on dbus BUS=system SIGNAL=NameOwnerChanged INTERFACE=org.freedesktop.DBus OBJPATH=/org/freedesktop/DBus ARG0=com.canonical.payments ARG2!=""

stop on desktop-end or \
    dbus BUS=system SIGNAL=NameOwnerChanged INTERFACE=org.freedesktop.DBus OBJPATH=/org/freedesktop/DBus ARG0=com.canonical.payments ARG2=""

respawn

script
    # XXX LP #1369692
    sleep 2

    exec /usr/bin/trust-stored-skeleton \
        --remote-agent DBusRemoteAgent --bus=session \
        --local-agent MirAgent \
        --trusted-mir-socket=/var/run/user/$(id -u)/mir_socket_trusted \
        --for-service UbuntuPayService \
        --with-text-domain ubuntu-pay-service \
        --store-bus session
end script
